# dev-embedded/gpbin/files/gphavagedstop
#!/sbin/openrc-run
# Copyright 1999-2016 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: $

user="root:root"

depend() {
    use alsasound
    before gpireteonline gpchronystop
    after local
}

_havastop() {
(
/etc/init.d/haveged stop
sleep 4
mount -a
sleep 6
mount -a
) &
}

start() {
	ebegin "Starting ${SVCNAME}"
        _havastop
        eend $?
}
