# dev-embedded/gpbin/files/gprtclock
#!/sbin/openrc-run
# Copyright 1999-2016 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: $

user="root:root"

depend() {
    use alsasound
    after bootmisc
}

_rtclock() {
echo 3072 >/sys/class/rtc/rtc0/max_user_freq
echo 3072 >/proc/sys/dev/hpet/max-user-freq
rc-service "${SVCNAME}" stop
}

start() {
	ebegin "Starting ${SVCNAME}"
        _rtclock
        eend $?
}
