# dev-embedded/gpbin/files/gpsetbub
#!/sbin/openrc-run
# Copyright 1999-2016 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: $

user="root:root"

depend() {
    use alsasound
    after netmount bubbleupnp
}


start() {
	ebegin "Starting ${SVCNAME}"
        /usr/bin/${SVCNAME}
        eend $?
}
