# dev-embedded/gpbin/files/ramdisk
#!/sbin/openrc-run

user="root:root"

depend() {
    need localmount
    after bootmisc
    before xdm
}

start() {
	ebegin "Starting ${SVCNAME}"
        /etc/default/${SVCNAME}.sh
        eend $?
}
